# monitor/gtk-monitor-scheme.path
[Unit]
Description=Sync GTK and caelestia schemes

[Path]
PathModified=%C/caelestia/scheme/current.txt
Unit=gtk-monitor-scheme.service

[Install]
WantedBy=default.target
